# deploy/ozon-dashboard.service
[Unit]
Description=Ozon Dashboard
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/ozon-dashboard
ExecStart=/opt/ozon-dashboard/.venv/bin/python dashboard_server.py 8765 127.0.0.1
Restart=always
RestartSec=5
User=ozon
Group=ozon
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target

